# dist/extras/cron.d/geoipupdate
# Run geoipupdate every week
# https://dev.maxmind.com/geoip/geoip2/geolite2/
@weekly root /usr/local/bin/geoipupdate
