# dist/extras/cron.d/update_firehol
# Run update_firehol.sh every day
# https://iplists.firehol.org/
@daily root /usr/local/bin/update_firehol.sh
